# dockerfiles/scrapyd-cron/crontabs/root
* */1 * * * curl http://scrapyd:6800/schedule.json -d project=example -d spider=example >> /var/log/cron.log 2>&1
